# source/luci-app-nfs/root/etc/init.d/nfs-luci
#!/bin/sh /etc/rc.common
# Unified NFS configuration wrapper for OpenWrt (UCI managed)
# Copyright (C) 2026 Gemini CLI

START=98
USE_PROCD=1

NFS_CONF=/etc/nfs.conf
EXPORTS_FILE=/etc/exports

write_share() {
        local cfg="$1"
        local path clients options
        config_get path "$cfg" path
        config_get clients "$cfg" clients "*"
        config_get options "$cfg" options "rw,sync,no_subtree_check"
        [ -n "$path" ] && echo "$path $clients($options)" >> "$EXPORTS_FILE"
}

start_service() {
        config_load nfs
        config_get_bool enabled globals enabled 0

        # If disabled in LuCI, ensure the native service is stopped
        if [ "$enabled" -eq 0 ]; then
                /etc/init.d/nfsd stop 2>/dev/null
                return
        fi

        logger -t "luci-app-nfs" -p notice "Generating NFS configurations from UCI"

        # Generate /etc/exports
        echo "# Generated from UCI /etc/config/nfs" > "$EXPORTS_FILE"
        config_foreach write_share share

        config_get_bool nfs_v3 globals nfs_v3 1
        config_get_bool nfs_v4 globals nfs_v4 1
        config_get threads globals threads 16
        config_get port globals port 2049
        config_get_bool udp globals udp 0
        config_get_bool tcp globals tcp 1

        # Generate /etc/nfs.conf
        cat > "$NFS_CONF" <<EOF
[nfsd]
threads=$threads
port=$port
udp=$( [ "$udp" -eq 1 ] && echo "y" || echo "n" )
tcp=$( [ "$tcp" -eq 1 ] && echo "y" || echo "n" )
vers3=$( [ "$nfs_v3" -eq 1 ] && echo "y" || echo "n" )
vers4=$( [ "$nfs_v4" -eq 1 ] && echo "y" || echo "n" )
EOF

        # Ensure the native daemon is enabled and start/reload it
        /etc/init.d/nfsd enable 2>/dev/null

        # Check if the native service is running; start if not, otherwise reload exports
        if /etc/init.d/nfsd status | grep -q "running"; then
                /usr/sbin/exportfs -r
        else
                /etc/init.d/nfsd start
        fi
}

stop_service() {
        logger -t "luci-app-nfs" -p notice "Stopping NFS service"
        /etc/init.d/nfsd stop 2>/dev/null
}

reload_service() {
        start_service
}

service_triggers() {
        procd_add_reload_trigger "nfs"
}
